# aidis-health.service
[Unit]
Description=AIDIS Health Check
After=aidis.service
Requires=aidis.service

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'curl -f http://localhost:8080/health || (systemctl restart aidis.service && logger -t aidis-health "Service unhealthy, restarted")'
StandardOutput=journal
StandardError=journal